# systemd/hermes-gc-responder-v2.service
[Unit]
Description=Hermes GC Responder v2 — @Ольга @Зинаида
After=network.target hermes-web-ui.service
Requires=hermes-web-ui.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /opt/zinaida/scripts/gc_responder_v2.py
WorkingDirectory=/opt/zinaida/scripts
Restart=always
RestartSec=10
User=root
StandardOutput=append:/var/log/gc-responder-v2.log
StandardError=append:/var/log/gc-responder-v2.log

[Install]
WantedBy=multi-user.target
